# slim-iso/airootfs/etc/skel/.config/autostart/vboxclient.desktop
[Desktop Entry]
Name=vboxclient
Comment=VirtualBox User Session Services
Exec=/usr/bin/VBoxClient-all

